# fs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/crond -c /etc/crontabs -f -l 5 -S