# fs/etc/services.d/openssh/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid "${USER_NAME:-openssh}" /usr/sbin/sshd -D -e -f /etc/openssh/sshd_config -p 2222
